# traefik/traefik.conf
[program:traefik]
command=/usr/local/bin/traefik.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/traefik.err.log
stdout_logfile=/var/log/traefik.out.log
