# example-site/example_site.conf
[program:example_site]
command=/usr/local/bin/example_site.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/example_site.err.log
stdout_logfile=/var/log/example_site.out.log
